# autostart/kbd-nocaps.desktop
[Desktop Entry]
Type=Application
Name=swap-caps
Exec=setxkbmap -option "ctrl:nocaps"
StartupNotify=false
Terminal=false
